# bin/99-second.monitor.rules
SUBSYSTEM=="drm", ACTION=="change", KERNEL=="card[0-9]*", RUN+="/usr/local/bin/monitors.sh"
